# rootdir/etc/systemd/user/pulseaudio.service.d/override.conf
[Service]
ExecStart=
EnvironmentFile=-/etc/sysconfig/pulseaudio
ExecStart=/usr/bin/pulseaudio --daemonize=no --log-target=journal $CONFIG
